# ios/runanywhere-sdks-latest/Playground/openclaw-hybrid-assistant/scripts/openclaw-voice.service
[Unit]
Description=OpenClaw Voice Assistant (Wake Word + STT + TTS)
After=openclaw-gateway.service
BindsTo=openclaw-gateway.service

[Service]
# Wait for gateway WebSocket (port 8082) to be ready
ExecStartPre=/bin/sleep 5
ExecStart=/home/runanywhere/runanywhere-sdks/Playground/openclaw-hybrid-assistant/build/openclaw-assistant --wakeword --url ws://localhost:8082
Restart=always
RestartSec=3
KillMode=process
Environment=HOME=/home/runanywhere
Environment="LD_LIBRARY_PATH=/home/runanywhere/runanywhere-sdks/sdk/runanywhere-commons/build/lib"

[Install]
WantedBy=default.target
